# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-hmc/meta-hmc/recipes-nvidia/fpga-ready-monitor/files/systemd/nvidia-fpga1-ready.service
[Unit]
Description=NVIDIA FPGA1 Ready service

[Service]
Restart=no
Type=simple
ExecStart=/bin/sh -c 'fpga1_ready.sh ready'
SyslogIdentifier=nvidia-fpga1-ready-service
