# meta-nvidia/meta-prime/meta-bmc-common/recipes-nvidia/set-hmc-time/files/set-hmc-time.timer
[Unit]
Description=Run set-hmc-timer service every hour 

[Timer]
OnUnitActiveSec=1h

#
# Upon periodic timer activation, systemd automatically activates set-hmc-time.service
#

[Install]
WantedBy=timers.target
