# meta-nvidia/meta-prime/meta-graceblackwell/meta-pcie-gpu/meta-bmc/recipes-phosphor/libmctp/files/systemd/mctp-usb-demux.socket
[Unit]
Description=MCTP USB 1-1.2.5 demux Unix domain socket

[Socket]
ListenSequentialPacket=@mctp-usb-mux-1-1-2-5
Accept=no

[Install]
WantedBy=sockets.target

