# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/powerstatus/files/check_cpu_boot_status.service
[Unit]
Description=NVIDIA CPU boots good checker service
Requires=nvidia-power-monitor.service
After=nvidia-power-monitor.service

[Service]
Type=simple
ExecStart=/usr/bin/env check_cpu_boot_status.sh
