# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-hmc/meta-hmc/recipes-phosphor/libmctp/files/systemd/fpga0-ap-recovery.target
[Unit]
Description=FPGA0 AP Recovery
PartOf= fpga0-ap-recovery.target