# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-phosphor/phosphor-gpio-monitor/phosphor-gpio-monitor/fan-fail@.service
[Unit]
Description=Fault service

Wants=obmc-chassis-hard-poweroff@0.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env error_log.sh fan-fail %I