# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-phosphor/phosphor-gpio-monitor/phosphor-gpio-monitor/power-fault@.service
[Unit]
Description=Fault service

[Service]
Type=oneshot
ExecStart=/usr/bin/env error_log.sh fault %I