# meta-nvidia/meta-prime/meta-bmc-common/recipes-nvidia/bmc-internal-network-config/files/ncsi-eth1-disable.service
[Unit]
Description= Disable NC-SI interface

Before=network.target

[Service]
Type=oneshot
ExecStart=ifconfig eth1 down

[Install]
WantedBy=multi-user.target
