# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-phosphor/libmctp/files/systemd/hmc-recovery.target
[Unit]
Description=HMC AP and ERoT Recovery
