# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-phosphor/libmctp/files/systemd/mctp-i2c8-ctrl.service
[Unit]
Description=MCTP I2C 8 control daemon
After=mctp-i2c8-demux.service obmc-power-on@0.target
Requires=mctp-i2c8-demux.service

[Service]
#Restart=always
RestartSec=10
ExecStart=/usr/bin/mctp-ctrl -m 1 -t 1 -f /usr/share/mctp/mctp_cfg_smbus8.json -n 8 -d 1 -v 1
SyslogIdentifier=mctp-i2c8-ctrl

[Install]
WantedBy=cpu-boot-done.service
