# meta-nvidia/meta-prime/meta-graceblackwell/meta-pcie-gpu/meta-bmc/recipes-phosphor/leds/files/power-led-controller.service
[Unit]
Description=Power LED Controller
After=mapper-wait@-xyz-openbmc_project-led-groups.service
After=mapper-wait@-xyz-openbmc_project-state-host0.service
Requires=xyz.openbmc_project.State.Boot.PostCode@0.service
After=xyz.openbmc_project.State.Boot.PostCode@0.service

[Service]
Restart=always
ExecStart=/usr/bin/env power-led-controller --config /usr/share/phosphor-led-manager/power-led-config.json

[Install]
WantedBy=multi-user.target
