# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/cpu-diag-status/files/cpu-diag-status.service
[Unit]
Description= CPU diagnostic status Service

[Service]
Type=oneshot
ExecStart=/usr/bin/env cpu-diag-status.sh 
SyslogIdentifier=cpu-diag-status

#
# This is also started periodically by cpu-diag-status.timer
#

[Install]
#WantedBy=multi-user.target
