# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/nvidia-code-mgmt/files/systemd/hmc-ready.service
[Unit]
Description=GB200NVL HMC Ready Service
Conflicts=hmc-recovery.target
StartLimitIntervalSec=1

[Service]
Restart=no
Type=simple
ExecStartPre=/bin/bash -c '/usr/bin/set-hmc-mux.sh FPGA; sleep 1'
ExecStart=/bin/echo "Stopping HMC AP Recovery Target"
SyslogIdentifier=gb200nvl-hmc-ready
