# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/nvidia-code-mgmt/files/systemd/hmc-notready.service
[Unit]
Description=GB200NVL HMC Not Ready Service
Requires=hmc-recovery.target
StartLimitIntervalSec=1

[Service]
Restart=no
Type=simple
ExecStart=/bin/echo "Starting HMC AP Recovery Target"
SyslogIdentifier=gb200nvl-hmc-notready