# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-phosphor/fans/phosphor-pid-control/fan-init.service
[Unit]
Description=Set Fan to Default Duty as Booting Up
Before=xyz.openbmc_project.fansensor.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/fan-init.sh

[Install]
WantedBy=sysinit.target