# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/rtc-detection/files/rtc-detection.service
[Unit]
Description=RTC detection Service
After=bmc-boot-complete.service

[Service]

RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/env rtc-detection.sh 
SyslogIdentifier=rtc-detection

[Install]
WantedBy=multi-user.target
