# meta-nvidia/meta-prime/recipes-nvidia/fpga-ready-monitor/files/nvidia-fpga-ready.target
[Unit]
Description=FPGA Ready
Conflicts=nvidia-fpga-notready.target
RefuseManualStop=yes

