# meta-nvidia/meta-prime/recipes-nvidia/fpga-ready-monitor/files/nvidia-fpga-notready.target
[Unit]
Description=FPGA Not Ready
Conflicts=nvidia-fpga-ready.target
RefuseManualStop=yes

