# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/powerctrl/files/nvidia-aux-power-force.service
[Unit]
Description=Aux cycle Host

[Service]
RemainAfterExit=yes
Type=oneshot
ExecStart=/usr/bin/env stbypowerctrl.sh aux_cycle force
SyslogIdentifier=auxcycle_force
