# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-bmc/recipes-nvidia/powerctrl/files/power-brake-enabled.target
[Unit]
Description=Power brake is currently enabled
After=
Conflicts=power-brake-disabled.target
RefuseManualStop=yes
