# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-hmc/recipes-nvidia/hmc-post-boot-cfg/files/bind_expanders.service
[Unit]
Description=Bind expanders afer fpga is up

[Service]
Type=simple
ExecStart=/usr/bin/bind_expanders.sh
SyslogIdentifier=bind_expanders

[Install]
WantedBy=nvidia-fpga-ready.target
