# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-hmc/meta-hmc/recipes-phosphor/phosphor-gpio-monitor/phosphor-gpio-monitor/sec-fpga-status@.service
[Unit]
Description=Detect and handle faults on FPGA_%i

[Service]
ExecStart=/usr/bin/env sec_fpga_fault_detect.sh %i
SyslogIdentifier=sec_fpga_fault_detect.sh
Type=oneshot
