# meta-nvidia/meta-prime/meta-gb200nvl/meta-bmc/recipes-nvidia/powerctrl/files/power-brake-disabled.target
[Unit]
Description=Power brake is currently disabled
After=
Conflicts=power-brake-enabled.target
RefuseManualStop=yes
