# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-hmc/recipes-phosphor/libmctp/files/systemd/mctp-spi0-demux.service
[Unit]
Description=MCTP SPI 0 demultiplexer daemon

[Service]
Type=notify
Restart=always
ExecStart=/usr/bin/mctp-demux-daemon astspi spi_config_file=/usr/share/mctp/mctp_cfg_spi0.json
SyslogIdentifier=mctp-spi0-demux
WatchdogSec=10