# meta-nvidia/recipes-phosphor/libmctp/libmctp/mctp-restart-notify.service
[Unit]
Description=MCTP Restart Notify service

[Service]
Type=oneshot

[Install]
WantedBy=reboot.target
