# meta-nvidia/meta-prime/meta-graceblackwell/meta-gb200nvl/meta-hmc/recipes-phosphor/libmctp/files/mctp-usb-ctrl.service
[Unit]
Description=MCTP USB control daemon

[Service]
Restart=always
Environment=MCTP_USB_CTRL_OPTS=null
EnvironmentFile=-/usr/share/mctp/mctp
ExecStart=/usr/bin/mctp-usb-ctrl $MCTP_USB_CTRL_OPTS
SyslogIdentifier=mctp-usb-ctrl

[Install]
WantedBy=nvidia-fpga-ready.target