# meta-nvidia/meta-prime/recipes-phosphor/console/obmc-console/systemd/obmc-console-ttyUSB5.socket
[Unit]
Description=Phosphor USB5 Console SSH Per-Connection socket

[Socket]
ListenStream=2203
Accept=yes

[Install]
WantedBy=sockets.target
