# meta-nvidia/meta-prime/recipes-phosphor/console/obmc-console/systemd/obmc-console-ttyUSB1.socket
[Unit]
Description=Phosphor USB1 Console SSH Per-Connection socket

[Socket]
ListenStream=2201
Accept=yes

[Install]
WantedBy=sockets.target
