# meta-nvidia/meta-prime/recipes-phosphor/console/obmc-console/systemd/obmc-console-ttyUSB4.socket
[Unit]
Description=Phosphor USB4 Console SSH Per-Connection socket

[Socket]
ListenStream=2202
Accept=yes

[Install]
WantedBy=sockets.target
