# meta-nvidia/meta-prime/meta-bmc-common/recipes-nvidia/sbiosbootaccess/files/cpu-boot-undone.service
[Unit]
Description=CPU boot undone service

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'control-bios-host-interface.sh boot-undone'
ExecStart=/usr/bin/env /usr/bin/delete-hi-user.sh
