# meta-nvidia/meta-prime/meta-bmc-common/recipes-nvidia/sbiosbootaccess/files/cpu-boot-done.service
[Unit]
Description=CPU boot done service

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'cpu-boot-done.sh'
